# rclone-to-kobo-watcher.path
[Path]
Unit=rclone-to-kobo-watcher.service
PathChanged=/home/gigi/Sync/Kobo

[Install]
WantedBy=multi-user.target
